# systemd/service-a.service
[Unit]
Description=Nginx Gateway Microservices — Service A (HTTP :3001)
After=network-online.target service-b.service service-c.service
Wants=network-online.target
Requires=service-b.service service-c.service

[Service]
Type=simple
User=microsvc
Group=microsvc
WorkingDirectory=/opt/nginx-gateway-microservices/services/service-a
Environment=PORT=3001
Environment=BIND_HOST=127.0.0.1
Environment=SERVICE_B_URL=http://service-b.internal:3002
Environment=NODE_ENV=production
ExecStartPre=/opt/nginx-gateway-microservices/scripts/wait-for-deps.sh
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=3s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=service-a

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
